# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command = /usr/sbin/service cron start
autostart=true
autorestart=false
startsecs=0
startretries=1
exitcodes=0
priority=1
